# archlive-install-server/releng/airootfs/etc/udev/rules.d/81-eth0-static-ip.rules
ACTION=="add", SUBSYSTEM=="net", ENV{INTERFACE}=="eth0", RUN+="/root/local/eth0-static-ip"
